# gamescope/galileo-mura-setup.service
[Unit]
Description=galileo mura setup
PartOf=graphical-session.target
After=graphical-session.target
ConditionFirmware=smbios-field(sys_vendor="Valve")
ConditionFirmware=smbios-field(product_name="Galileo")

[Service]
Type=oneshot
ExecStart=/usr/bin/galileo-mura-setup
EnvironmentFile=%t/gamescope-environment
